# demoApp.service
[Unit]
Description=myapp
After=syslog.target

[Service]
ExecStart=java -jar /opt/demoApp/demo.jar

[Install]
WantedBy=multi-user.target
